# supervisord.ini
[program:vulwarning]
command=/opt/sss/vulwarning/vulwarning
directory=/opt/sss/vulwarning
autostart=true
autorestart=true
startsecs=10
startretries=3
priority=1
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/opt/sss/vulwarning/stdout.log
